# camera-streamer.service
[Unit]
Description=Camera Streamer Service
After=network.target

[Service]
ExecStart=/opt/camera-streamer/scripts/stream_camera.sh
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

# Environment variables
Environment=CAMERA_DEVICE=/dev/video0
Environment=RTSP_PORT=8554
Environment=STREAM_NAME=camera

[Install]
WantedBy=multi-user.target
