# mediamtx.service
[Unit]
Wants=network.target
[Service]
ExecStart=/usr/local/bin/mediamtx
[Install]
WantedBy=multi-user.target
